# docker/supervisord.conf
[supervisord]
nodaemon=true

# Twistd is the web service that runs the flask app
[program:vss365todaytwitter]
command=python3 ./listener.py
directory=/app
